# services/pacbackup.service
[Unit]
Description=Back up installed packages list

[Service]
Type=simple
ExecStart=/bin/sh -c '/usr/bin/pacman -Qqe > /etc/packages.txt'
