# services/pacmonitor.service
[Unit]
Description=Check for updates

[Service]
Type=simple
ExecStart=/usr/bin/pacman -Sy

[Install]
WantedBy=multi-user.target
