# services/snebu.timer
[Unit]
Description=Backup weekly

[Timer]
OnCalendar=weekly
Unit=snebu.service

[Install]
WantedBy=multi-user.target
