# services/system/pacmonitor.service
[Unit]
Description=Run pacman -Sy

[Service]
Type=simple
ExecStart=/usr/local/bin/pacmonitor
